# config/crontab
0 8 * * * php /var/www/html/home/scripts/server/cronjob.php